# provision/wifi/setup-wifi.service
[Unit]
Description=Configure Wi-Fi from /boot/wifi_credentials.txt
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-wifi.sh
RemainAfterExit=yes

[Install]
# This hooks the service into the multi-user boot process.
WantedBy=multi-user.target

